# config/crontab
* * * * * /usr/bin/python3 /opt/scripts/script.py >> /var/log/cron.log 2>&1

